# ansible/files/server/crontab
# MajiTask_v1 Production Crontab
# Automated maintenance tasks for MajiTask_v1 server
# Updated: June 28, 2025

# SSL Certificate renewal (Let's Encrypt) - Daily at 5 PM
0 17 * * * /usr/bin/certbot renew --quiet --deploy-hook "systemctl reload nginx"

# MajiTask_v1 application backup - Daily at 2 AM
0 2 * * * /root/bin/server_backup.sh >> /var/log/majitask_v1_backup.log 2>&1

# Database optimization - Weekly on Sunday at 3 AM
0 3 * * 0 /usr/bin/mysqlcheck -h127.0.0.1 -umajitask_v1 -pmajitask_secure_2024! --optimize majitask_v1 >> /var/log/db_maintenance.log 2>&1

# Log rotation and cleanup - Daily at 4 AM
0 4 * * * /usr/sbin/logrotate /etc/logrotate.conf

# System updates check - Weekly on Monday at 6 AM
0 6 * * 1 /usr/bin/apt update && /usr/bin/apt list --upgradable > /var/log/system_updates.log 2>&1

# Application health check - Every 5 minutes
*/5 * * * * /bin/systemctl is-active --quiet majitask_v1 || (/bin/systemctl restart majitask_v1 && echo "$(date): MajiTask_v1 restarted" >> /var/log/majitask_v1_health.log)

# Cleanup temporary files - Daily at midnight
0 0 * * * /usr/bin/find /tmp -type f -atime +7 -delete
0 0 * * * /usr/bin/find /home/apps/majitask_v1/majitask/uploads/temp -type f -atime +1 -delete 2>/dev/null || true

# Monitor disk space - Every 6 hours
0 */6 * * * /bin/df -h | /usr/bin/awk '$5 > 85 {print "WARNING: Disk usage high on " $1 ": " $5}' | /usr/bin/mail -s "Disk Space Alert" majitask.fun@gmail.com

# Application log monitoring - Every hour (check for errors)
0 * * * * /usr/bin/tail -n 100 /home/apps/majitask_v1/majitask/logs/error.log | /bin/grep -i "error\|critical\|fatal" | /usr/bin/mail -s "MajiTask_v1 Error Alert" majitask.fun@gmail.com 2>/dev/null || true
