# majitask/majitask.service
[Unit]
Description=MajiTask_v1 - Task Management Application
Documentation=file:///home/apps/majitask_v1/VM_DEPLOYMENT_GUIDE.md
After=network.target
Wants=network.target

[Service]
Type=simple
User=majitask_v1
Group=majitask_v1
WorkingDirectory=/home/apps/majitask_v1/majitask
Environment=NODE_ENV=production
Environment=PORT=3863
Environment=HOST=0.0.0.0
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=majitask_v1

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/apps/majitask_v1/majitask

[Install]
WantedBy=multi-user.target
